# junkyard/swayosd-wayfire/99-swayosd-backlight.rules
# Allow members of the "video" group to change display backlight + keyboard
# LED brightness. Required by swayosd-client --brightness on systems where
# /sys/class/backlight/*/brightness is otherwise root-only.
ACTION=="add", SUBSYSTEM=="backlight", RUN+="/bin/chgrp video /sys%p/brightness", RUN+="/bin/chmod g+w /sys%p/brightness"
ACTION=="add", SUBSYSTEM=="leds", KERNEL=="*::kbd_backlight", RUN+="/bin/chgrp video /sys%p/brightness", RUN+="/bin/chmod g+w /sys%p/brightness"
